# playbooks/roles/monitor/files/otel-collector.service
[Unit]
Description=OpenTelemetry Collector
After=network.target

[Service]
ExecStart=/usr/local/bin/otelcol --config=/etc/otel-collector/config.yaml
Restart=always

[Install]
WantedBy=multi-user.target
